# state/prod/application/elastic/kibana/files/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
directory=/opt/kibana/bin/
user=root
autostart = true
autorestart = true
startsecs = 5
