# state/test/city/ly-txds-3p-3j/python/supervisor/files/ly_log.conf
[program:ly_log]
command=/opt/logstash/bin/logstash -f /etc/logstash/ly_log.json
directory=/opt/logstash/bin
user=root
autostart = true
autorestart = true
startsecs = 5
environment=PATH="/usr/local/services/jdk/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin",JAVA_HOME="/usr/local/services/jdk"
