# state/prod/city/ly-txds-3p-3j/python/supervisor/files/udplog_bj.conf
[program:udplog_bj]
command=/opt/logstash/bin/logstash -f /etc/logstash/udplog_bj.json
directory=/opt/logstash/bin
user=root
autostart = true
autorestart = true
startsecs = 5
environment=PATH="/usr/local/services/jdk/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin",JAVA_HOME="/usr/local/services/jdk"
